# data/systemd/snapd.target
[Unit]
Description=Snapd fully seeded target
Requires=basic.target
After=basic.target
AllowIsolate=yes
